# KitchenSyncD_mk2_desktop/other/arpping.service
[Unit]
Description=ARP Fixer to ensure reliable kitchensync

[Service]
ExecStart=/usr/local/bin/arpping_svc.sh

[Install]
WantedBy=multi-user.target
