# KitchenSyncD_mk2_desktop/other/paranoid.service
[Unit]
Description=Paranoia mac changer

[Service]
ExecStart=/usr/local/bin/paranoid_svc.sh

[Install]
WantedBy=multi-user.target
